# systemd/user/battery-alert.service
[Unit]
Description=Battery alerts for Hyprland

[Service]
ExecStart=%h/.config/hypr/scripts/battery_alert.sh

